# auto_add_music.service
[Unit]
Description=Auto Add Music - watch %h/Music/AutoAdd and organize into library
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/carlos/scripts/music_rel/auto_add_music.py
WorkingDirectory=/home/carlos
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target